# quadlet/worb.container
# ==============================================================================
# worb App - Quadlet Configuration (Production)
# Generated by CodeB CLI v3.0.4
# ==============================================================================
# Environment: production
# Port: 4000
# Network: host (Storage server direct connection)
# Database: db.codeb.kr (Storage server)
# WebSocket: ws.codeb.kr (Centrifugo)
# ==============================================================================

[Unit]
Description=worb CMS Production Container
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/codeb-kr/worb:latest
ContainerName=worb
AutoUpdate=registry

# Network: host mode for Storage server direct connection
Network=host

# Basic environment variables
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=HOSTNAME=0.0.0.0

# Storage server connection (db.codeb.kr - 64.176.226.119)
Environment=DATABASE_URL=postgresql://worb_user:worb_secure_pass@db.codeb.kr:5432/worb?schema=public
Environment=REDIS_URL=redis://db.codeb.kr:6379/0
Environment=REDIS_PREFIX=worb:

# Centrifugo connection (ws.codeb.kr - 141.164.42.213)
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api
Environment=CENTRIFUGO_API_KEY=pRMupNs6HlGp7G6xkPsAFrI8hN4g6U0G
Environment=CENTRIFUGO_SECRET=of0KuRFjjzhq5LlBURCuKqzTUAA08hwL

# Additional environment variables (with secrets) - create separately on server
EnvironmentFile=/etc/containers/systemd/worb.env

# Health Check
HealthCmd=curl -sf http://localhost:4000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# Resource limits
PodmanArgs=--memory=2g --cpus=2

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
