# quadlet/worb-staging.container
# ==============================================================================
# worb App - Quadlet Configuration (Staging)
# Generated by CodeB CLI v3.0.4
# ==============================================================================
# Environment: staging
# Port: 3002
# Network: host (Storage server direct connection)
# Database: db.codeb.kr (Storage server, staging DB)
# WebSocket: ws.codeb.kr (Centrifugo)
# ==============================================================================

[Unit]
Description=worb CMS Staging Container
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/codeb-kr/worb:staging
ContainerName=worb-staging
AutoUpdate=registry

# Network: host mode for Storage server direct connection
Network=host

# Basic environment variables
Environment=NODE_ENV=staging
Environment=PORT=3002
Environment=HOSTNAME=0.0.0.0

# Storage server connection (db.codeb.kr - 64.176.226.119) - Staging DB
Environment=DATABASE_URL=postgresql://worb_user:worb_secure_pass@db.codeb.kr:5432/worb_staging?schema=public
Environment=REDIS_URL=redis://db.codeb.kr:6379/1
Environment=REDIS_PREFIX=worb-staging:

# Centrifugo connection (ws.codeb.kr - 141.164.42.213)
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api
Environment=CENTRIFUGO_API_KEY=pRMupNs6HlGp7G6xkPsAFrI8hN4g6U0G
Environment=CENTRIFUGO_SECRET=of0KuRFjjzhq5LlBURCuKqzTUAA08hwL

# Additional environment variables (with secrets) - create separately on server
EnvironmentFile=/etc/containers/systemd/worb-staging.env

# Health Check
HealthCmd=curl -sf http://localhost:3002/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# Resource limits (Staging has less than Production)
PodmanArgs=--memory=1g --cpus=1

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
