# quadlet/workb-cms-staging.container
[Unit]
Description=WorkB CMS Staging Container
After=network-online.target
Wants=network-online.target
After=workb-cms-postgres.service
After=workb-cms-redis.service
Requires=workb-cms-postgres.service
Requires=workb-cms-redis.service

[Container]
ContainerName=workb-cms-staging
Image=localhost/workb-cms:staging

# Port mapping (Staging uses different port)
PublishPort=3201:3000

# Environment
Environment=NODE_ENV=staging
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0

# Environment file with secrets
EnvironmentFile=/etc/containers/systemd/workb-cms-staging.env

# Health check
HealthCmd=curl -sf http://localhost:3000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# Resource limits (Staging - lower limits)
PodmanArgs=--memory=1g --memory-swap=1.5g --cpus=1 --pids-limit=128

# Network
Network=workb-network

# DNS (container discovery)
AddHost=workb-cms-postgres:10.89.1.10
AddHost=workb-cms-redis:10.89.1.11

# Logging
LogDriver=journald

# Security
NoNewPrivileges=true
ReadOnlyRootfs=false

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# Graceful shutdown
ExecStop=/usr/bin/podman stop -t 30 %n

[Install]
WantedBy=multi-user.target default.target
